# ids.service
[Unit]
Description=Python Intrusion Detection System
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/nafnlaus/ids/base-ids.py
Restart=on-failure
RestartSec=5
User=nafnlaus
StandardOutput=append:/var/log/ids.log
StandardError=append:/var/log/ids_error.log

[Install]
WantedBy=multi-user.target

